# config/com.aios.data-collect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aios.data-collect</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kamallamgade/Documents/AIOS Setup/AI-Work-Assistant-System/.venv/bin/python</string>
        <string>/Users/kamallamgade/Documents/AIOS Setup/AI-Work-Assistant-System/scripts/collect.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/kamallamgade/Documents/AIOS Setup/AI-Work-Assistant-System/data/collect.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kamallamgade/Documents/AIOS Setup/AI-Work-Assistant-System/data/collect.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/kamallamgade/Documents/AIOS Setup/AI-Work-Assistant-System</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
